# myflaskapp.service
[Unit]
Description=My Flask App
After=network.target

[Service]
User=yourusername
Group=www-data
WorkingDirectory=/path/to/your/app
Environment="PATH=/path/to/your/venv/bin"
ExecStart=/usr/local/bin/docker-compose up

[Install]
WantedBy=multi-user.target
